# configs/releng/airootfs/root/customize/customize/autostart/org.gnome.Console.desktop
[Desktop Entry]
Type=Application
Name=Console
Exec=kgx
